# node2/systemd/system/pg-failover.timer
[Timer]
OnCalendar=*-*-* 22:05:00
Persistent=true
Unit=pg-failover-restart.service

[Install]
WantedBy=timers.target